# time-synchronized.target
[Unit]
Description=Kernel Time Synchronized
RefuseManualStart=yes
